# deployment/systemd/apachas-backup-monitor.timer
[Unit]
Description=Check A Pachas encrypted backup freshness hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10m
AccuracySec=5m
Persistent=true
Unit=apachas-backup-monitor.service

[Install]
WantedBy=timers.target
