# deployment/systemd/apachas-backup.service
[Unit]
Description=A Pachas encrypted data backup
Documentation=https://github.com/adpablos/apachas/blob/main/docs/deployment.md
After=docker.service
Wants=docker.service
RequiresMountsFor=/var/lib/docker /var/backups/apachas

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=-/etc/apachas/backup.env
ExecStart=/opt/apachas/scripts/backup_data.sh
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ReadOnlyPaths=/var/lib/docker/volumes/apachas_api-data/_data
ReadOnlyPaths=/etc/apachas
ReadWritePaths=/var/backups/apachas
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
LockPersonality=true
MemoryDenyWriteExecute=true
